# etc/cron.d/braincase
# M H D M Y user command

*/1 * * * * root braincase-user-sync
*/5 * * * * root braincase-dropbox-auto-setup
0 0 * * *  root braincase-backups
*/30 * * * * root braincase-memorize
