# etc/cron.d/braincase
# M H D M Y user command

*/1 * * * * root braincase-user-sync
0 0 * * *  root braincase-backups --run