# etc/cron.d/braincase
# M H D M Y user command

*/1 * * * * root braincase-user-sync
